# www-servers/jetty/files/jetty
#!/sbin/runscript

JETTY_PID=/var/run/jetty.pid

start() {
	ebegin "Starting Jetty"
	start-stop-daemon --start --quiet \
		--chdir /usr/share/jetty \
		--user jetty:jetty --pidfile ${JETTY_PID} \
		--make-pidfile --background \
		--exec java -- -jar start.jar --ini=/etc/jetty/start.ini --daemon

	eend $? || return 1
}

stop() {
	ebegin "Stopping Jetty"
	start-stop-daemon --stop \
		--pidfile ${JETTY_PID}
}

# kate: replace-tabs false; hl: Bash
