# www-apps/tabbyapi/files/tabbyapi.socket
[Socket]
Accept=false
ListenStream=5000
Service=tabbyapi-proxy.service

[Install]
WantedBy=sockets.target
