# net-analyzer/ensure-speed/files/ensure-speed.timer
[Unit]
Description=Timed tester for bandwidth speed

[Timer]
OnCalendar=*-*-* 0/12:00:00

[Install]
WantedBy=multi-user.target
