# net-analyzer/ensure-speed/files/ensure-speed.service
[Unit]
Description=Speed tester

[Service]
Type=oneshot
ExecStart=/usr/bin/ensure-speed -s -d 50 -u 10 -p 15 -l /var/log/ensure-speed.log
