# sys-devel/prelink/files/prelink.service
[Unit]
Description=Prelink

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-prelink
