# dev-db/dolt/files/dolt.service
[Unit]
Description=Dolt SQL server
Documentation=https://docs.dolthub.com/
After=network.target

[Service]
Type=simple
User=dolt
Group=dolt
WorkingDirectory=/var/lib/dolt
# Use 127.0.0.1 so Dolt does not enable the default Unix socket under /tmp/mysql.sock.
ExecStart=/usr/bin/dolt sql-server --host 127.0.0.1
Restart=on-failure
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
